# config/crontab
# 每分钟执行备份，flock防止重复执行
* * * * * flock -n /tmp/rclone.lock rclone sync mywebdav:nas/我的文件/ /data --log-file=/var/log/rclone-sync.log --log-level=INFO --stats=10m